# build/package/debian/control/postinst
#!/bin/sh
# postinst script

set -e

NAME=butler-xbee-gateway

case "$1" in
	configure)
	;;

	upgrade)
	;;

	abort-upgrade|abort-remove|abort-deconfigure)
	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

if [ -x "/etc/init.d/${NAME}" ]; then
	update-rc.d ${NAME} defaults >/dev/null
	invoke-rc.d ${NAME} start || exit $?
fi

exit 0
